# root/etc/s6-overlay/s6-rc.d/init-mod-vaapi-amdgpu/run
#!/usr/bin/with-contenv bash

echo "**** Setting up AMD VAAPI drivers ****"

PLEX_TRANSCODER="/usr/lib/plexmediaserver/Plex Transcoder"
PLEX_TRANSCODER_ORIG="/usr/lib/plexmediaserver/Plex Transcoder.orig"

# ============================================================
# Step 1: Create wrapper for Plex Transcoder
# This ensures our libva/Mesa stack is used instead of Plex's
# ============================================================

# Only create wrapper if we haven't already
if [ ! -f "$PLEX_TRANSCODER_ORIG" ]; then
    echo "Creating Plex Transcoder wrapper..."
    
    # Move original transcoder
    mv "$PLEX_TRANSCODER" "$PLEX_TRANSCODER_ORIG"
    
    # Create wrapper script that injects our library path
    cat > "$PLEX_TRANSCODER" << 'EOF'
#!/bin/bash
# AMD VAAPI wrapper - ensures modern Mesa/libva is used
export LD_LIBRARY_PATH="/vaapi-amdgpu/lib:${LD_LIBRARY_PATH}"
export LIBVA_DRIVERS_PATH="/vaapi-amdgpu/lib/dri"
export LIBVA_DRIVER_NAME="radeonsi"

# Fix shader cache permission issue - use Plex's config dir
export MESA_SHADER_CACHE_DIR="/config/Library/Application Support/Plex Media Server/Cache/mesa-shader-cache"
mkdir -p "$MESA_SHADER_CACHE_DIR" 2>/dev/null

exec "/usr/lib/plexmediaserver/Plex Transcoder.orig" "$@"
EOF
    
    chmod +x "$PLEX_TRANSCODER"
    echo "Wrapper created successfully"
else
    echo "Transcoder wrapper already exists"
fi

# ============================================================
# Step 2: Set up symlinks in Plex's cache (belt & suspenders)
# ============================================================

PLEX_VA_CACHE="/config/Library/Application Support/Plex Media Server/Cache/va-dri-linux-x86_64"
mkdir -p "$PLEX_VA_CACHE"

# Remove any old/stale drivers Plex may have downloaded
rm -f "$PLEX_VA_CACHE"/*.so* 2>/dev/null

# Symlink our modern drivers
for driver in /vaapi-amdgpu/lib/dri/*.so; do
    if [ -f "$driver" ]; then
        ln -sf "$driver" "$PLEX_VA_CACHE/$(basename "$driver")"
        echo "Linked driver: $(basename "$driver")"
    fi
done

# ============================================================
# Step 3: Verify setup
# ============================================================

if [ -f /usr/share/libdrm/amdgpu.ids ]; then
    echo "amdgpu.ids present"
else
    echo "WARNING: amdgpu.ids not found"
fi

echo "**** AMD VAAPI setup complete ****"

