# root/etc/s6-overlay/s6-rc.d/init-mod-vaapi-amdgpu/up
/etc/s6-overlay/s6-rc.d/init-mod-vaapi-amdgpu/run

